# bin/ensembleiot.service
[Unit]
Description=PiOT Framework
Before=systemd-user-sessions.service
After=network-online.target
Wants=network-online.target

[Service]
TimeoutStartSec=0

ExecStart=BINPATH/run.sh
Type=oneshot
RemainAfterExit=yes
User=USER

[Install]
WantedBy=multi-user.target
