# bin/ensembleiot.service
[Unit]
Description=Ensemble IoT Daemon
Before=systemd-user-sessions.service
After=network.target
Requires=network.target

[Service]
TimeoutStartSec=0

ExecStart=BINPATH/run.sh
Type=simple
RemainAfterExit=no
User=USER

ExecStop=
User=USER

[Install]
WantedBy=multi-user.target
